# files/clown-mount-ephemeral.conf
start on stopped mountall
exec /usr/local/bin/clown-mount-ephemeral