# files/clown-init-ephemeral.conf
start on starting mountall
exec /usr/local/bin/clown-init-ephemeral
